# app/mouse_remote_control/R/app/22-vx-uibc.rules

ACTION!="add", GOTO="uibc_exit"
KERNEL=="hidraw2", SUBSYSTEM=="hidraw", RUN+="/etc/udev/startup-uibc.sh"
LABEL="uibc_exit"

